# root/etc/services.d/skydns/run
#!/usr/bin/execlineb -P


skydns -machines="http://127.0.0.1:4001,http://127.0.0.1:4001" -addr="0.0.0.0:53" -domain="devoply.local." -verbose=true -nameservers="8.8.8.8:53,8.8.4.4:53"
